# _delugeTools/raspberry/40-deluge-midi.rules
ACTION=="add", SUBSYSTEM=="block", ATTR{partition}=="*", RUN+="/usr/local/bin/DelugeMidiUdev.sh"

